# roles/adminscripts/templates/cron.d/node_down_reminders.j2
28 6 5,19 * 1 root /usr/local/bin/ring-admin send downreminders >/dev/null
