# roles/adminscripts/templates/cron.d/node_down_deactivation.j2
28 6 27 * * root /usr/local/bin/ring-admin purge machines
