# roles/adminscripts/templates/cron.d/ssh_keyscan.j2
37 22 * * * root /usr/local/bin/ring-admin scan hostkeys 1>&2
