# roles/ring_auth/templates/cron.d/auth-copy-keys-to-opt-keys.j2
*/15 * * * * root /usr/local/bin/auth-copy-keys-to-opt-keys
