# roles/ring_auth/templates/cron.d/auth_sync_graphite_credentials.j2
*/15 * * * * root /usr/local/bin/auth_sync_graphite_credentials
