# roles/graphite/templates/cron.d/import_graphite_users.j2
12 * * * * graphite ssh auth.infra.ring.nlnog.net cat /opt/graphite/users.list | /usr/local/bin/ring-graphite import users
