# roles/apt/files/apt-daily-upgrade.timer
[Timer]
OnCalendar=
OnCalendar=daily
RandomizedDelaySec=23h
