# roles/ntp/templates/cron.d/ntp.j2
@reboot root /usr/sbin/ntpdate {{ ntp_default_server }}
