# roles/ntp/templates/cron.d/ntp.j2
@reboot root /bin/sleep 300 && /usr/sbin/ntpdate {{ ntp_default_server }} >/dev/null
