# roles/scamper/templates/scamper.conf
# Scamper deamon
start on runlevel [2345]
stop on shutdown

respawn
respawn limit 10 5

exec /usr/bin/scamper -P 23456 -M $(hostname)
