# roles/inframail/templates/cron.d/ring-aptsummary
01 10 * * * apt /usr/local/bin/ring-aptsummary
