# roles/ring_sqa/templates/systemd-ring-sqa.service.j2
[Unit]
Description=RING SQA (IPv{{ "6" if afi == "ipv6" else "4" }}

[Service]
ExecStart=/usr/local/bin/ring-sqad {{ "-6" if afi == "ipv6" else "" }}

